# ubus_controller/files/ubus_controller.init
#!/bin/sh /etc/rc.common

USE_PROCD=1
START=50

start_service() {
    local enabled

    procd_open_instance
        procd_set_param command /usr/bin/ubus_controller
        procd_close_instance
}

service_triggers() {
	procd_add_reload_trigger "ubus_controller"
}

reload_service() {
	stop
	start
}